# Arch-silverblue/iso/airootfs/etc/systemd/system/greetd.service.d/live.conf
# Run greetd with the live-installer kiosk config (the stock
# /etc/greetd/config.toml belongs to the greetd package).
# Inert unless the live ISO ships greetd (LIVE_GUI_SESSION=1 in distro.conf).
[Service]
ExecStart=
ExecStart=/usr/bin/greetd --config /etc/greetd/live.toml
